# etc/hotplug.d/iface/99_fix_dns
[ "$ACTION" = "ifdown" -a "$INTERFACE" = "vpn" ] && /www/bin/flush_dns_fix.sh
